# xyz.openbmc_project.GARPControl.service.in
[Unit]
Description=Phosphor GARP Control
After=systemd-networkd.service xyz.openbmc_project.Network.service
[Service]
ExecStart=/usr/bin/env phosphor-garpcontrol
SyslogIdentifier=phosphor-garpcontrol

[Install]
WantedBy=@SYSTEMD_TARGET@
